# packaging/app-software-updates.cron
2 3 * * * root /usr/sbin/software-updates
