# packaging/app-software-updates-cache.cron
30 2 * * * root /usr/bin/yum --noplugins --enablerepo=* clean metadata >/dev/null 2>&1
